# zarf/systemd.service
[Unit]
Description=X1 Carbon Gen9 Thermal limiter

[Service]
Type=simple
ExecStart=/usr/sbin/x1-cpu-manager.sh

[Install]
WantedBy=multi-user.target
